# dev/ks-system-home/DEBIAN/prerm
#!/bin/bash

. /usr/share/debconf/confmodule

if [ "$1" != 'upgrade' ]; then
    systemctl disable --now syncthing@kael.timer
    rm "$HOME/.config/syncthing/config.xml"
fi

