# pkg-src/deb/insync-portable/DEBIAN/prerm
#!/bin/bash

if [ "$1" != 'upgrade' ]; then
    systemctl disable --now insync-portable@
    rm -Rf /opt/insync*
    if [ "$1" == "purge" ]; then
        rm -Rf /home/*/.config/Insync
    fi
fi

